# HW6/HW6.service
[Unit]
Description=HW6

[Service]
ExecStart=/home/tms/DevOps_TMS/HW6/task1.sh

[Install]
WantedBy=multi-user.target
